# scripts/systemd/dropins/v2x-perception-link-health.service.d/v2x-alerts.conf
# Escalate link-health failures to the V2X alert notifier and send a single
# recovery notice when the check next passes. Installed as a drop-in so the
# live unit file is never overwritten.
[Unit]
OnFailure=v2x-alert@%n.service
OnSuccess=v2x-alert-recovered@%n.service
