# scripts/systemd/v2x-cloudflared-drive.service
[Unit]
Description=V2X Cloudflare Drive Quick Tunnel
Wants=network-online.target v2x-drive.service
After=network-online.target v2x-drive.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=path
Group=path
WorkingDirectory=/home/path
EnvironmentFile=-/etc/v2x-drive-tunnel.env
ExecStart=/home/path/V2XCarla/v2x-backend/scripts/launch-cloudflared-drive-tunnel.sh
Restart=always
RestartSec=10
TimeoutStopSec=30
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
