# scripts/systemd/v2x-feed-recovery.service
[Unit]
Description=Run post-recovery validation when perception feeds return
OnFailure=v2x-alert@%n.service

[Service]
Type=oneshot
User=path
Group=path
WorkingDirectory=/home/path/V2XCarla/v2x-backend
StateDirectory=v2x-feed-recovery
StateDirectoryMode=0700
Environment=STATE_DIR=/var/lib/v2x-feed-recovery
ExecStart=/home/path/V2XCarla/v2x-backend/scripts/check-feed-recovery.sh
