# scripts/systemd/v2x-perception.service
[Unit]
Description=V2X Perception Pipeline
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=path
Group=path
WorkingDirectory=/home/path/V2XCarla/v2x-backend
Environment=AWS_REGION=us-west-2
Environment=AWS_PROFILE=Path-Emerging-Dev-147229569658
Environment=PYTHONUNBUFFERED=1
Environment=PYTHON_BIN=/home/path/V2XCarla/perception-venv/bin/python
Environment=TORCH_HOME=/home/path/.cache/torch
Environment=V2X_PERCEPTION_MODEL_PATH=/home/path/V2XCarla/v2x-backend/apps/perception/yolov8n.pt
Environment=V2X_PERCEPTION_VIDEO_PATHS=v2x-backend-cam-ch1,v2x-backend-cam-ch2,v2x-backend-cam-ch3,v2x-backend-cam-ch4
Environment=V2X_VIDEO_SESSION_API_BASE_URL=https://w0j9m7dgpg.execute-api.us-west-1.amazonaws.com
Environment=V2X_PERCEPTION_CONFIDENCE=0.4
Environment=V2X_PERCEPTION_SHOW_LIVE=false
Environment=V2X_PERCEPTION_STREAM_HOST=0.0.0.0
Environment=V2X_PERCEPTION_STREAM_PORT=8090
Environment=V2X_PERCEPTION_UPLOAD_MIN_INTERVAL_SEC=1.0
Environment=V2X_PERCEPTION_STALE_SECONDS=15
Environment=V2X_PERCEPTION_INFERENCE_STALE_SECONDS=10
Environment=V2X_PERCEPTION_RECONNECT_INITIAL_SEC=1
Environment=V2X_PERCEPTION_RECONNECT_MAX_SEC=30
Environment=V2X_PERCEPTION_OPEN_TIMEOUT_MS=10000
Environment=V2X_PERCEPTION_READ_TIMEOUT_MS=10000
Environment=V2X_PERCEPTION_LIVE_HLS_FRAGMENTS=4
Environment=V2X_PERCEPTION_CAPTURE_HLS_FRAGMENTS=1
Environment=V2X_PERCEPTION_CLOCK_HLS_FRAGMENTS=5
Environment=V2X_PERCEPTION_CLOCK_GLITCH_GRACE_SEC=2
Environment=V2X_PERCEPTION_PROACTIVE_RENEW_SEC=240
Environment=V2X_PERCEPTION_TERMINAL_READ_FAILOVER_SEC=5
Environment=V2X_PERCEPTION_MEDIA_CLOCK_MIN_LATENCY_MS=-1000
Environment=V2X_PERCEPTION_MEDIA_CLOCK_MAX_LATENCY_MS=10000
Environment=V2X_PERCEPTION_CAPTURE_BACKEND=ffmpeg_nvdec
Environment=V2X_PERCEPTION_INFERENCE_WORKERS=2
Environment=V2X_PERCEPTION_CROSS_CAMERA_VEHICLE_ASSOCIATION=false
Environment=V2X_PERCEPTION_FFMPEG_BIN=/usr/bin/ffmpeg
Environment="OPENCV_FFMPEG_CAPTURE_OPTIONS=live_start_index;-1"
Environment=V2X_DETECTIONS_ENDPOINT=https://w0j9m7dgpg.execute-api.us-west-1.amazonaws.com/detections
EnvironmentFile=-/etc/v2x-perception.env
ExecStartPre=/usr/bin/test -x /home/path/V2XCarla/perception-venv/bin/python
ExecStartPre=/usr/bin/test -r /home/path/V2XCarla/v2x-backend/apps/perception/yolov8n.pt
ExecStartPre=/usr/bin/test -r /home/path/.cache/torch/hub/checkpoints/mobilenet_v3_small-047dcff4.pth
ExecStartPre=/usr/bin/test -x /usr/bin/ffmpeg
ExecStart=/home/path/V2XCarla/v2x-backend/scripts/launch-perception.sh
Restart=always
RestartSec=10
TimeoutStopSec=60
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
