# scripts/systemd/dropins/v2x-perception.service.d/v2x-alerts.conf
# Alert when the perception service itself enters a failed state. No
# OnSuccess here: the service is long-running and only stops on operator
# action, so success notices would be noise.
[Unit]
OnFailure=v2x-alert@%n.service
